# contrib/systemd/homelab-verify.timer
# homelab-verify.timer — fire the gate check nightly.
[Unit]
Description=Run homelab-iac gate verification nightly

[Timer]
# 03:30 local, with jitter so it's not on a hard boundary.
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=900
Persistent=true          # catch up if the box was asleep at 03:30

[Install]
WantedBy=timers.target
