# contrib/systemd/homelab-verify.service
# homelab-verify.service — one-shot nightly gate check.
# Install on whatever box holds the repo + vault password (your MBP if it's
# always-on, or a small always-on LAN host). NOT on the PVE host.
#
#   sudo cp contrib/systemd/homelab-verify.* /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now homelab-verify.timer
#
# Runs as your user so it uses your SSH agent/keys and the repo checkout.
# Edit User= and the paths to match your environment.
[Unit]
Description=homelab-iac nightly gate verification
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=straderb
Environment=HOMELAB_IAC_REPO=/opt/homelab-iac
# Point at your ansible if it's in a venv, e.g.:
# Environment=PATH=/opt/homelab-iac/.venv/bin:/usr/bin:/bin
ExecStart=/opt/homelab-iac/contrib/bin/nightly-verify.sh
# Don't hang forever if a VM is wedged:
TimeoutStartSec=600
