# contrib/systemd/homelab-verify@.service
# Instantiated nightly verification unit. The instance is the unprivileged
# Linux account that owns the checkout and its SSH credentials:
#   systemctl start homelab-verify@alice.service
[Unit]
Description=homelab-iac non-disruptive verification for %i
Wants=network-online.target
After=network-online.target
ConditionOS=linux
ConditionPathIsDirectory=/opt/homelab-iac/.venv

[Service]
Type=oneshot
User=%i
WorkingDirectory=/opt/homelab-iac
ExecStart=/opt/homelab-iac/.venv/bin/ansible-playbook verify.yml --vault-password-file .vault_pass --extra-vars notify_on_success=false
TimeoutStartSec=600
UMask=0077
