# roles/svc_download/files/proxies/ll-proxy.socket
# /etc/systemd/system/ll-proxy.socket
[Unit]
Description=lazylibrarian UI proxy socket (LAN -> vpn netns veth)
# FreeBind lets the socket bind before 192.168.1.31 is configured at boot.
Wants=network-online.target
After=network-online.target

[Socket]
ListenStream=192.168.1.31:5299
FreeBind=yes

[Install]
WantedBy=sockets.target
