# roles/svc_download/files/proxies/ll-proxy.service
# /etc/systemd/system/ll-proxy.service  (activated by ll-proxy.socket)
[Unit]
Description=LazyLibrarian UI proxy (LAN -> 10.77.0.2:5299)
Requires=vpn-netns.service dl-lazylibrarian.service
After=vpn-netns.service dl-lazylibrarian.service

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --exit-idle-time=5min 10.77.0.2:5299
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
