# roles/svc_download/files/proxies/prowlarr-proxy.service
# /etc/systemd/system/prowlarr-proxy.service  (activated by prowlarr-proxy.socket)
[Unit]
Description=Prowlarr UI proxy (VLAN -> 10.77.0.2:9696)
Requires=vpn-netns.service dl-prowlarr.service
After=vpn-netns.service dl-prowlarr.service

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --exit-idle-time=5min 10.77.0.2:9696
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
