# roles/svc_download/files/leak-canary.service
# /etc/systemd/system/leak-canary.service
[Unit]
Description=VPN jail leak canary
After=vpn-netns.service

[Service]
Type=oneshot
EnvironmentFile=/etc/homelab-notify.env
ExecStart=/usr/local/sbin/leak-canary.sh
