# roles/svc_media/files/backup-media.service
# /etc/systemd/system/backup-media.service
[Unit]
Description=Nightly svc-media backup (RomM DB dump + appdata tars) to NFS
RequiresMountsFor=/srv/backups

[Service]
Type=oneshot
User=homelab
Group=homelab
Environment=XDG_RUNTIME_DIR=/run/user/10001
ExecStart=/usr/local/sbin/backup-media.sh
