# roles/svc_infra/templates/beszel-agent.container.j2
{{ ansible_managed | comment }}
# ~homelab/.config/containers/systemd/beszel-agent.container   (rootless)
# Reports this host's metrics to the beszel-hub catalog entry. Outbound-only
# (connects TO the hub) — no published port, no firewalld rule needed here.
[Unit]
Description=Beszel monitoring agent

[Container]
Image={{ beszel_agent_image }}
ContainerName=beszel-agent
# Host networking: the agent monitors the HOST, so it must see the host's real
# NICs (a bridge/pasta netns would report the container's eth0 instead). It is
# also REQUIRED here specifically — the hub runs on THIS host, and a rootless
# bridge/pasta container cannot reach its own host's published LAN IP
# (connection refused); the host netns reaches {{ hostvars[infra_host].ansible_host }}:8090 normally.
Network=host
Environment=HUB_URL=http://{{ hostvars[infra_host].ansible_host }}:8090
Environment=TOKEN={{ vault_beszel_token }}
# Quoted: the hub public key is "ssh-ed25519 AAAA..." (contains a space), and
# an unquoted Environment= splits on whitespace — which would truncate KEY to
# "ssh-ed25519" and turn the base64 into a second, bogus bare env var.
Environment="KEY={{ vault_beszel_key }}"
Volume=%h/appdata/beszel-agent:/var/lib/beszel-agent:Z

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
