# roles/svc_infra/templates/infra-app.container.j2
{{ ansible_managed | comment }}
# ~homelab/.config/containers/systemd/{{ item.key }}.container   (rootless)
# RENDERED from roles/svc_infra/templates/infra-app.container.j2, driven by
# the infra_apps catalog entry "{{ item.key }}". Published on svc-infra's
# LAN IP, NOT loopback — Caddy runs on {{ media_host }}, a different host,
# and needs to reach this directly (see caddy_services in
# group_vars/all/main.yml, https://{{ item.key }}.{{ service_domain }});
# svc_infra's files.yml opens ui_port to lan_cidr via firewalld to match.
[Unit]
Description={{ item.key }}

[Container]
Image={{ item.value.image }}
ContainerName={{ item.key }}
{% if item.value.network_mode | default('') == 'host' %}
PodmanArgs=--network=host
{% else %}
PublishPort={{ ansible_host }}:{{ item.value.ui_port }}:{{ item.value.container_port | default(item.value.ui_port) }}
{% for p in item.value.extra_ports | default([]) %}
PublishPort={{ ansible_host }}:{{ p.port }}:{{ p.port }}/{{ p.protocol }}
{% endfor %}
{% endif %}
{% for volume in item.value.volumes %}
Volume={{ volume }}
{% endfor %}
Environment=TZ={{ timezone }}
{% for env_key, env_value in (item.value.env | default({})).items() %}
Environment="{{ env_key }}={{ env_value }}"
{% endfor %}

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
