# roles/svc_infra/templates/netbox-db.container.j2
{{ ansible_managed | comment }}
[Unit]
Description=NetBox PostgreSQL

[Container]
Image={{ netbox_images.postgres }}
ContainerName=netbox-db
Network=netbox.network
# postgres:18+ stores data in a major-version subdir and wants the mount at
# /var/lib/postgresql (NOT .../data — that path errors out on 18+). See
# https://github.com/docker-library/postgres/pull/1259
Volume=%h/appdata/netbox-db:/var/lib/postgresql:Z
Environment=TZ={{ timezone }}
Environment=POSTGRES_DB=netbox
Environment=POSTGRES_USER=netbox
Environment=POSTGRES_PASSWORD={{ vault_netbox_db_password }}
HealthCmd=pg_isready -q -d netbox -U netbox
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=20s
Notify=healthy

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
