# roles/svc_infra/templates/paperless-redis.container.j2
{{ ansible_managed | comment }}
[Unit]
Description=Paperless-ngx Redis broker

[Container]
Image={{ paperless_images.redis }}
ContainerName=paperless-redis
Network=paperless.network
# valkey (redis-wire-compatible) broker — see paperless_images in apps.yml
# for why NOT the stock redis image (its non-root user can't write this
# homelab-owned bind mount under rootless podman → MISCONF wedges the broker).
# Its own dir, deliberately OUTSIDE appdata/paperless/, so the nightly
# `tar paperless` backup doesn't sweep it — just a cache/task-queue broker,
# not worth backing up.
Volume=%h/appdata/paperless-redis:/data:Z
Environment=TZ={{ timezone }}

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
