# bundle-reference/svc-download/proxies/hydra-proxy.service
# /etc/systemd/system/hydra-proxy.service  (activated by hydra-proxy.socket)
[Unit]
Description=NZBHydra2 UI proxy (VLAN -> 10.77.0.2:5076)
Requires=vpn-netns.service dl-nzbhydra2.service
After=vpn-netns.service dl-nzbhydra2.service

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --exit-idle-time=5min 10.77.0.2:5076
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
