# bundle-reference/svc-download/proxies/sab-proxy.socket
# /etc/systemd/system/sab-proxy.socket
#
# Exposes SAB's UI on the VM's VLAN address without giving the jail any
# route to the LAN. systemd-socket-proxyd runs in the init ns and connects
# to 10.77.0.2 over the /30 — replies return over the connected route, so
# the jail never needs (or gets) a path to 192.168.0.0/16.
# Socket activation also means the proxy only runs while someone is using it.
[Unit]
Description=SABnzbd UI proxy socket (VLAN -> vpn netns veth)

[Socket]
ListenStream=192.168.1.31:8080
BindIPv6Only=ipv6-only

[Install]
WantedBy=sockets.target
