# bundle-reference/svc-download/proxies/sonarr-proxy.service
# /etc/systemd/system/sonarr-proxy.service  (activated by sonarr-proxy.socket)
[Unit]
Description=Sonarr UI proxy (VLAN -> 10.77.0.2:8989)
Requires=vpn-netns.service dl-sonarr.service
After=vpn-netns.service dl-sonarr.service

[Service]
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --exit-idle-time=5min 10.77.0.2:8989
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=yes
