# roles/svc_download/templates/beszel-agent.container.j2
{{ ansible_managed | comment }}
# /etc/containers/systemd/beszel-agent.container   (rootful, init ns)
# NOT catalog-managed (no download_apps entry, deliberately not in the VPN
# jail — it needs normal LAN egress to reach the hub, which host-backstop.nft
# only allows once vault_beszel_token/_key are configured).
[Unit]
Description=Beszel monitoring agent

[Container]
Image={{ beszel_agent_image }}
ContainerName=beszel-agent
# Host networking: the agent monitors the HOST, so it must see the host's real
# NICs (a bridge/pasta netns would report the container's eth0 instead). On the
# init-ns host this also means egress to the hub follows the host backstop's
# 8090 allow rule directly.
Network=host
Environment=HUB_URL=http://{{ hostvars[infra_host].ansible_host }}:8090
Environment=TOKEN={{ vault_beszel_token }}
# Quoted: the hub public key is "ssh-ed25519 AAAA..." (contains a space), and
# an unquoted Environment= splits on whitespace — which would truncate KEY to
# "ssh-ed25519" and turn the base64 into a second, bogus bare env var.
Environment="KEY={{ vault_beszel_key }}"
Volume=/srv/appdata/beszel-agent:/var/lib/beszel-agent:Z

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
