# bundle-reference/svc-media/quadlets/jellyfin.container
# ~homelab/.config/containers/systemd/jellyfin.container   (rootless)
#
# Ports: 8096 published on all interfaces (not loopback) — LAN clients
# (girlfriend's Chromecast w/ the Jellyfin Android TV app, phones at home)
# direct-play against http://192.168.1.30:8096, while tailnet clients go
# through Caddy at https://svc-media.<tailnet>.ts.net:8444. Scope 8096 in
# firewalld to the LAN zone only.
#
# Transcoding: the 3700X host has NO iGPU — this is CPU (software) transcode.
# Fine for 1-2 concurrent 1080p transcodes on 4 vCPU; prefer direct-play
# clients, and see docs/v2-video.md for the GPU-passthrough upgrade path.
# /transcodes on tmpfs keeps that churn off the VM disk.
[Unit]
Description=Jellyfin

[Container]
Image=docker.io/jellyfin/jellyfin:latest
ContainerName=jellyfin
Label=io.containers.autoupdate=registry
PublishPort=8096:8096
Volume=%h/appdata/jellyfin/config:/config:Z
Volume=%h/appdata/jellyfin/cache:/cache:Z
Volume=/srv/media/tv:/media/tv:ro
Volume=/srv/media/movies:/media/movies:ro
Tmpfs=/transcodes:size=4g
Environment=JELLYFIN_PublishedServerUrl=http://192.168.1.30:8096

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
