# roles/svc_media/files/quadlets/audiobookshelf.container
# svc-media — ROOTLESS Quadlet, place in:
#   ~homelab/.config/containers/systemd/audiobookshelf.container
# as the 'homelab' user (uid/gid 10001), with lingering enabled:
#   loginctl enable-linger homelab
# Rootless is fine here (no netns tricks needed on this VM); the container
# runs as 10001 which matches TrueNAS dataset ownership, so no PUID env and
# no userns-vs-NFS mapping problem.
[Unit]
Description=Audiobookshelf

[Container]
Image=ghcr.io/advplyr/audiobookshelf:latest
ContainerName=audiobookshelf
PublishPort=127.0.0.1:13378:80
Volume=%h/appdata/audiobookshelf/config:/config:Z
Volume=%h/appdata/audiobookshelf/metadata:/metadata:Z
Volume=/srv/media/audiobooks:/audiobooks
Volume=/srv/media/ebooks:/ebooks:ro
# v5: ABS built-in scheduled backups (Settings > Backups) default to
# /config/backups — now bind-mounted to the NFS backup dataset, so the
# UI default Just Works. No :Z on NFS mounts.
Volume=/srv/backups/audiobookshelf:/config/backups

[Service]
ExecStartPre=/bin/sh -c 'stat /srv/media/. >/dev/null && findmnt -t nfs4 -M /srv/media >/dev/null && stat /srv/backups/. >/dev/null && findmnt -t nfs4 -M /srv/backups >/dev/null'
Restart=on-failure

[Install]
WantedBy=default.target
