# roles/svc_media/files/quadlets/ntfy.container
# ~homelab/.config/containers/systemd/ntfy.container   (rootless)
#
# The alert sink for the whole homelab: deploy pings, the svc-download leak
# canary, and the disk-space alarms all POST here (ntfy_url =
# http://192.168.1.30:8080). Published on the LAN IP so those LAN/loopback
# clients reach it directly; also reverse-proxied by Caddy at ntfy.<domain>
# for the browser UI. Anonymous read+write on the LAN — set vault_ntfy_token
# + auth in server.yml if you want to lock it down.
[Unit]
Description=ntfy notification server

[Container]
Image=docker.io/binwiederhier/ntfy:latest
ContainerName=ntfy
Exec=serve
PublishPort=192.168.1.30:8080:80
Volume=%h/appdata/ntfy/cache:/var/cache/ntfy:Z
Volume=%h/appdata/ntfy/etc:/etc/ntfy:Z
Environment=TZ=America/New_York
Environment=NTFY_BASE_URL=http://192.168.1.30:8080

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
