# roles/svc_media/files/quadlets/seerr.container
# ~homelab/.config/containers/systemd/seerr.container   (rootless)
#
# Media request front-end — Seerr, the unified successor to Overseerr +
# Jellyseerr (supports Jellyfin/Plex/Emby in one app). This is what you hand
# other people to request shows/movies; it drives Sonarr/Radarr at
# 192.168.1.31:8989/7878 (the socket proxies on the download VM, admitted from
# LAN_ADMIN). Reverse-proxied by Caddy at seerr.<domain>. Same port (5055) as
# the old Jellyseerr, so no proxy/backstop/firewall changes.
#
# Config is Overseerr/Jellyseerr-compatible — to carry over an existing setup,
# see the migration note in docs/media.md.
[Unit]
Description=Seerr

[Container]
Image=ghcr.io/seerr-team/seerr:latest
ContainerName=seerr
PublishPort=127.0.0.1:5055:5055
Volume=%h/appdata/seerr:/app/config:Z
Environment=TZ=America/New_York
Environment=PORT=5055
# The official image ships no init process — add one for zombie reaping.
PodmanArgs=--init

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
