# roles/mon/files/homelab-diskalert.timer
[Unit]
Description=Run the disk-space alarm every 15 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
