# roles/mon/files/homelab-diskalert.service
[Unit]
Description=Disk-space ntfy alarm
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/disk-alert.sh
